# profile.d/environment.sh
export PAGER=less
#if [ $TERM = "screen.linux" ] ; then
	#export TERM=linux
#fi
